# internal/service/example.service
[Unit]
Description=gometrum - metrics to Home Assistant via MQTT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gometrum --config /etc/gometrum.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
